# etc/systemd/mochad-check.timer
[Unit]
Description=Monitor Mochad and restart if needed

[Timer]
OnBootSec=30
OnUnitActiveSec=240

[Install]
WantedBy=multi-user.target
